# solr_daemon
root@Workstation-oyj-X555QG:/etc/init.d# cat solr
#!/bin/sh
### BEGIN INIT INFO
# Provides:          solr_daemon
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Manage solr_daemon server
# Description:       Start, stop, restart solr server for a specific application.
### END INIT INFO
solr_start(){
     su - oyj -c "/home/oyj/rails_root/myblog/solr_start.sh"
}
  

solr_stop(){
     SOLR_PID=$(ps -ef | grep solr | grep "port=8983" | awk '{ print $2}')
     if [ ! $SOLR_PID ]]
     then
       echo "Solr is not running"
     else
          kill -9 $SOLR_PID
     fi
}
#cd /home/oyj/rails_root/myblog
#su - oyj -c "/home/oyj/rails_root/myblog/solr_start.sh"

case "$1" in
   start)
     echo "Starting solr daemon"
     solr_start
   ;;

   stop)
     echo  "Stopping solr daemon"
     solr_stop
   ;;
esac
